# modules/supervisor/templates/supervisorTask.conf.erb
[program:<%= @name %>]
user=<%= @user %>
command=<%= @command %>
autorestart=<%= @autorestart %>
numprocs=<%= @numprocs %>
process_name=%(program_name)s_%(process_num)s